# ops/aliyun/medical-credit-api.service.example
[Unit]
Description=medical-credit-assessment Aliyun API proxy
After=network.target

[Service]
WorkingDirectory=/var/www/medical-credit-api/current
ExecStart=/usr/bin/node aliyun-api/server.js
Restart=always
RestartSec=3
EnvironmentFile=/var/www/medical-credit-api/.env

[Install]
WantedBy=multi-user.target
